# config/supervisor/supervisord.conf
[supervisord]
nodaemon = true

[program:rsyslog]
command = /usr/sbin/rsyslogd -n
autorestart = false
priority = 2000

[program:crond]
command = /usr/sbin/crond -n
priority = 1000
autorestart = true

[program:opendmarc]
command = /usr/sbin/opendmarc -c /etc/opendmarc.conf -P /var/run/opendmarc/opendmarc.pid
autorestart = false
priority = 500

[program:postfix]
process_name = master
command = /opt/postfix
directory = /etc/postfix
autorestart = false
priority = 200
